# gunicorn.service
[Unit]
Description=gunicorn daemon for it-desk-tool
Requires=gunicorn.socket
After=network.target

[Service]
User=viki
Group=www-data
WorkingDirectory=/home/viki/it-desk-tool/backend

# --- THE FINAL FIX ---
# We now force the production settings directly into the command using the --env flag.
# This makes it impossible for the service to use the wrong settings.
ExecStart=/home/viki/it-desk-tool/backend/venv/bin/gunicorn \
          --env DJANGO_SETTINGS_MODULE=helpdesk.settings.production \
          --workers 3 \
          --bind unix:/run/gunicorn.sock \
          helpdesk.wsgi:application

[Install]
WantedBy=multi-user.target