# gunicorn.service
# E:\it-admin-tool\gunicorn.service

[Unit]
Description=gunicorn daemon for it-desk-tool
Requires=gunicorn.socket
After=network.target

[Service]
User=viki
Group=www-data
WorkingDirectory=/home/viki/it-desk-tool/backend

# --- THE FINAL, GUARANTEED FIX ---
# This now uses the configuration file, which is the most robust method.
ExecStart=/home/viki/it-desk-tool/backend/venv/bin/gunicorn --config gunicorn_config.py helpdesk.wsgi:application

[Install]
WantedBy=multi-user.target